# mmoinsweeper.service.j2
[Unit]
Description=MMOInsweeper Game Server
After=network.target

[Service]
ExecStart=/usr/local/bin/mmoinsweeper --port 8080
Restart=always
User=mmoinsweeper
Group=mmoinsweeper
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
